# files/system/usr/lib/systemd/user/flatpak-portal.service.d/50-nocblue-clean-activation-env.conf
[Service]
# flatpak-portal launches helper sandboxes for apps such as Steam Runtime.
# Host allocator preload does not belong in those child sandbox environments.
UnsetEnvironment=LD_PRELOAD GTK_USE_PORTAL
