# files/system/usr/lib/systemd/user/dbus-.service.d/50-nocblue-clean-activation-env.conf
[Service]
# Native D-Bus services can be activated by portals. Forcing them back through
# portals can deadlock the request path, notably FileManager1 -> Nautilus.
UnsetEnvironment=GTK_USE_PORTAL
