# files/system/usr/lib/systemd/user/xdg-desktop-portal-gnome.service.d/50-nocblue-clean-activation-env.conf
[Service]
# Portal services must not force their own requests back through portals or
# leak host allocator preload into helper sandboxes.
UnsetEnvironment=LD_PRELOAD GTK_USE_PORTAL
