# files/system/usr/lib/systemd/system/brew-upgrade.service.d/50-nocblue-brew-proxy.conf
[Unit]
ConditionPathIsSymbolicLink=
ConditionFileIsExecutable=/home/linuxbrew/.linuxbrew/bin/brew
