# files/system/usr/lib/systemd/user/nix-portable-gc.timer
[Unit]
Description=Weekly Nix store cleanup

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=2h

[Install]
WantedBy=timers.target
