# files/system/usr/lib/systemd/user/nocblue-flatpak-icon-fixes.service
[Unit]
Description=Fix known Flatpak SVG icons for Qt renderers
After=graphical-session.target nocblue-user-defaults.service

[Service]
Type=oneshot
ExecStart=/usr/bin/nocblue-fix-flatpak-icons

[Install]
WantedBy=default.target
