# files/system/usr/lib/systemd/user/nocblue-flatpak-setup.timer
[Unit]
Description=Install nocblue user Flatpaks
ConditionUser=!@system

[Timer]
OnStartupSec=45s

[Install]
WantedBy=timers.target
