# files/system/usr/lib/systemd/user/nocblue-adwsteamgtk.timer
[Unit]
Description=Retry nocblue Steam GTK skin setup
ConditionPathExists=!%h/.local/state/nocblue/adwsteamgtk.done

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
