# files/system/usr/lib/systemd/user/nocblue-icon-theme-sync.service
[Unit]
Description=Sync nocblue Candy icon theme for sandboxed apps

[Service]
Type=oneshot
ExecStart=/usr/bin/nocblue-sync-icon-theme
