# files/system/usr/lib/systemd/user/nocblue-game-shortcuts-sync.service
[Unit]
Description=Sync Steam and Lutris shortcuts for launcher search
After=nocblue-flatpak-setup.service
ConditionUser=!@system
ConditionPathIsDirectory=%h

[Service]
Type=oneshot
ExecStart=/usr/bin/nocblue-sync-game-shortcuts
