# files/system/usr/lib/systemd/user/nocblue-adwsteamgtk.service
[Unit]
Description=Apply the nocblue Steam GTK skin
After=network-online.target graphical-session.target nocblue-user-defaults.service nocblue-noctalia-theme.service
Wants=network-online.target nocblue-noctalia-theme.service
ConditionPathExists=!%h/.local/state/nocblue/adwsteamgtk.v3.done

[Service]
Type=oneshot
ExecStart=/usr/bin/nocblue-apply-adwsteamgtk
