# files/system/usr/lib/systemd/user/noctalia.service.d/50-nocblue-standard-malloc.conf
[Service]
# Noctalia crashes during plugin installation with hardened_malloc.
# Retain secureblue's libno_rlimit_as preload without the allocator.
Environment=LD_PRELOAD=libno_rlimit_as.so
