# files/system/usr/lib/systemd/system/rpm-ostreed-automatic.service.d/zz-nocblue-verify-update.conf
[Service]
# secureblue's automatic-update pre-check only recognizes ghcr.io/secureblue/*.
# Replace it for downstream nocblue images and keep signed-image verification.
ExecStartPre=
ExecStartPre=/usr/bin/nocblue-verify-update-signature
