# files/system/usr/lib/systemd/user/xdg-desktop-portal.service.d/50-nocblue-clean-activation-env.conf
[Service]
# Portal services must not force their own requests back through portals.
UnsetEnvironment=GTK_USE_PORTAL
